# scripts/launchd/com.nanika.nen-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: install.sh replaces __HOME__ with the user's actual $HOME path. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanika.nen-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.alluka/bin/nen-daemon</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.alluka/bin:__HOME__/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Redirect stdout/stderr to log files under ~/.alluka/logs/ -->
    <key>StandardOutPath</key>
    <string>__HOME__/.alluka/logs/nen-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.alluka/logs/nen-daemon.log</string>

    <!-- launchd restarts the daemon on any exit (including SIGTERM from `nen-daemon stop`) -->
    <key>KeepAlive</key>
    <true/>

    <!-- Minimum seconds between restarts; prevents tight crash loops -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- Start when the agent is bootstrapped / user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Working directory must exist before start -->
    <key>WorkingDirectory</key>
    <string>__HOME__/.alluka</string>
</dict>
</plist>
